# cups/rootfs/etc/services.d/cups/run
#!/usr/bin/with-contenv bashio
cupsd -f
