# cups/rootfs/etc/services.d/cups/run
#!/usr/bin/with-contenv bashio
mkdir -p /var/run/cups/certs
mkdir -p /var/log/cups
chmod 0755 /var/run/cups/certs
echo "Starting CUPS daemon..."
cupsd -f -c /etc/cups/cupsd.conf -l /var/log/cups/cups.log
