# cups/rootfs/etc/services.d/cups/run
#!/usr/bin/with-contenv bashio
chmod 0755 /var/run/cups/certs
cupsd -f -c /etc/cups/cupsd.conf
